# launchd/com.os1.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.os1.app</string>

    <!--
      Launch the proper .app bundle via `open -W` so macOS registers it as
      a GUI app (window can be drawn, Aqua session context, AppleScript
      reachable). Launching the raw Mach-O binary directly leaves it
      headless — process runs but no window ever appears.

      -W waits for the app to quit, so launchd's KeepAlive policy applies.
      -n forbids opening a second copy if one is already running.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/open</string>
        <string>-W</string>
        <string>-n</string>
        <string>/Users/gaganarora/Desktop/my projects/hermes-desktop-os1/dist/OS1.app</string>
    </array>

    <!--
      No secrets here. OS1 reads:
        ELEVENLABS_API_KEY  ← Keychain "io.elevenlabs.api-key"
        ELEVENLABS_AGENT_ID ← Keychain "io.elevenlabs.agent-id"
        ORGO_API_KEY        ← Keychain "ai.orgo.mac.api-key"
    -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>OS1_REALTIME_ORGO_TOOLSETS</key>
        <string>core,screen,files,admin</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/Users/gaganarora/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/gaganarora</string>
    </dict>

    <!--
      Only restart on crash. If the user clicks the red X, leave it closed.
    -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/gaganarora/.os1/os1.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gaganarora/.os1/os1.err</string>
</dict>
</plist>
