# launchd/com.os1.coo.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.os1.coo</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>__OS1_HOME__/.os1/coo/daemon.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__OS1_PATH__</string>
        <key>HOME</key>
        <string>__OS1_HOME__</string>
        <key>OS1_COO_INTERVAL_SECONDS</key>
        <string>1800</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>__OS1_HOME__/.os1/coo/daemon.log</string>
    <key>StandardErrorPath</key>
    <string>__OS1_HOME__/.os1/coo/daemon.err</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
